# .config/misc/crontab
*/5 * * * * DISPLAY=:0 /usr/local/bin/mailsync

